# gateway/com.deepgem.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.deepgem.gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jarvis/jarvis-config/gateway/start-gateway.sh</string>
    </array>

    <!-- Keep alive with throttle to prevent restart loops -->
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/jarvis/Library/Logs/Jarvis/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarvis/Library/Logs/Jarvis/gateway-error.log</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/jarvis/jarvis-config/gateway</string>
</dict>
</plist>
